# package/DEBIAN/postinst
#!/bin/bash
# post-installation script

# Enable and start the vulnreport daemon
systemctl enable vulnreport.service
systemctl start vulnreport.service