# package/DEBIAN/prerm
#!/bin/bash
# pre-removal script

# Stop and disable vulnreport daemon
systemctl stop vulnreport.service
systemctl disable vulnreport.service