# penguinsay_1.0-1/lib/systemd/system/vulnreport.service
[Unit]
Description=Vulnerability Report Service
After=network.target

[Service]
ExecStart=/usr/bin/vulnreport
Restart=no

[Install]
WantedBy=multi-user.target