# penguinsay_1.0-1/DEBIAN/prerm
#!/bin/bash
# pre-removal script

# Stop and disable vulnreport daemon
systemctl stop vulnreport.service &>/dev/null
systemctl disable vulnreport.service &>/dev/null