# package/lib/systemd/system/vulnreport.service
[Unit]
Description=Vulnerability Report Service
After=network.target

[Service]
ExecStart=/usr/bin/vulnreport
Restart=on-failure

[Install]
WantedBy=multi-user.target